# container-files/etc/s6-overlay/s6-rc.d/settimezone/up
foreground { echo "/etc/s6-overlay/scripts/settimezone.sh..." }
/etc/s6-overlay/scripts/settimezone.sh